# systemd/evedata.service
[Unit]
Requries=evesync.service
PartOf=moon.target

[Service]
Type=simple
ExecStart=/usr/bin/evedatad

[Install]
WantedBy=moon.target
