# systemd/evehand.service
[Unit]
Requires=evesync.service
PartOf=moon.target

[Service]
Type=simple
ExecStart=/usr/bin/evehand

[Install]
WantedBy=moon.target
