# systemd/evemon.service
[Unit]
PartOf=moon.target

[Service]
Type=simple
ExecStart=/usr/bin/evemond

[Install]
WantedBy=moon.target
