# systemd/moon.target
[Unit]
Description=Sysmoon base target
After=multi-user.target

[Install]
WantedBy=multi-user.target
