# systemd/syssync.service
[Unit]
Requires=sysmoon.service
PartOf=moon.target

[Service]
Type=simple
ExecStart=/usr/bin/syssyncd

[Install]
WantedBy=moon.target
